# deploy/ansible/templates/daemon.conf
description "daemon for site"

start on startup
stop on shutdown
respawn

env NODE_ENV={{env}}
env PORT={{site_port}}

script
  /usr/bin/npm start --prefix={{site_path}}
end script
